# deploy/systemd/gateway.service
[Unit]
Description=The-Block Gateway
After=network.target

[Service]
User=block
ExecStart=/usr/local/bin/gateway --listen 0.0.0.0:9000
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
