# deploy/systemd/node.service
[Unit]
Description=The-Block Node
After=network.target

[Service]
User=block
EnvironmentFile=/etc/the-block/tls/node.env
Environment=TLS_MANIFEST_PATH=/etc/the-block/tls/node/tls-manifest.json
ExecStartPre=/usr/bin/env bash -c '/usr/local/bin/contract tls stage --input /etc/the-block/tls --name node --service node:required=/etc/the-block/tls/node --env-file /etc/the-block/tls/node.env --force'
ExecReload=/usr/bin/env bash -c '/usr/local/bin/contract tls stage --input /etc/the-block/tls --name node --service node:required=/etc/the-block/tls/node --env-file /etc/the-block/tls/node.env --force && systemctl kill -s HUP %n'
ExecStart=/usr/local/bin/node --rpc 0.0.0.0:8545 --p2p 0.0.0.0:8001
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
